# systemd/uni-sync.service
[Unit]
Description=Hourly University Data Sync
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=michael
ExecStart=/bin/bash /home/michael/cachyos-setup/uni-sync.sh
StandardOutput=append:/home/michael/cachyos-setup/cron_debug.log
StandardError=append:/home/michael/cachyos-setup/cron_debug.log

[Install]
WantedBy=multi-user.target
