# systemd/uni-sync.timer
[Unit]
Description=Run University Sync every hour

[Timer]
OnCalendar=hourly
Persistent=true
Unit=uni-sync.service

[Install]
WantedBy=timers.target
